# template.service
[Unit]
Description=Network Service

[Service]
ExecStart=ping 8.8.8.8